# macos/launchd/com.tds.lmde.sync-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tds.lmde.sync-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/CHANGE_ME/workplace/tds-utils/bin/launch-sync-monitor</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/lmde-sync-monitor.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/lmde-sync-monitor.err.log</string>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
